# amazonia-sensor.service
[Unit]
Description=Sensor Amazonia IoT – Simulador MQTT
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%RUNTIME_USER%
WorkingDirectory=%RUNTIME_DIR%
ExecStart=%RUNTIME_NODE% %RUNTIME_DIR%/sensor.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=-%RUNTIME_DIR%/.env

[Install]
WantedBy=multi-user.target